# infra/systemd/controlpanel-agent.service
[Unit]
Description=ControlPanel Node Agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/controlpanel-agent --config /etc/controlpanel/agent.yaml
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target

